# deploy/capybara-api.service
[Unit]
Description=Capybara Runner API Server
After=network.target

[Service]
Type=simple
User=cc2
WorkingDirectory=/home/cc2/game01/server
Environment=DATABASE_URL=sqlite:data.db
Environment=HOST=127.0.0.1
Environment=PORT=3000
Environment=CORS_ORIGINS=http://localhost
ExecStart=/home/cc2/game01/server/target/release/capybara-server
Restart=on-failure
RestartSec=5

# Security hardening
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
ReadWritePaths=/home/cc2/game01/server /home/cc2/backups/capybara

[Install]
WantedBy=multi-user.target
